# deploy/raspberry-pi/nodos-backend-restart.service
[Unit]
Description=Reinicio programado de nodos-backend

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart nodos-backend
