# deploy/raspberry-pi/nodos-db-backup.service
[Unit]
Description=Backup diario de nodos.db al pendrive
After=local-fs.target
Wants=network.target

[Service]
Type=oneshot
User=__USER__
Environment=APP_DIR=__APP_DIR__
Environment=BACKUP_MOUNT=/mnt/nodos-backup
ExecStart=__APP_DIR__/deploy/raspberry-pi/backup-db.sh
