# deploy/raspberry-pi/nodos-db-backup.timer
[Unit]
Description=Backup diario de nodos.db

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
